# packaging/container/rootfs/etc/s6-overlay/s6-rc.d/blackhole-threats/run
#!/bin/sh
exec 2>&1

timestamp() {
    date -u '+%Y-%m-%dT%H:%M:%SZ'
}

escape_value() {
    printf '%s' "$1" | sed 's/\\/\\\\/g; s/"/\\"/g'
}

log_line() {
    level="$1"
    message="$2"
    shift 2

    printf 'time=%s level=%s msg="%s"' "$(timestamp)" "$level" "$(escape_value "$message")"
    while [ "$#" -gt 0 ]; do
        key="${1%%=*}"
        value="${1#*=}"
        printf ' %s="%s"' "$key" "$(escape_value "$value")"
        shift
    done
    printf '\n'
}

log_info() {
    log_line info "$@"
}

log_error() {
    log_line error "$@"
}

DEFAULT_CONF="/config/blackhole-threats.yaml"
CONF_PATH="${BLACKHOLE_THREATS_CONF:-$DEFAULT_CONF}"
EXTRA_OPTS="${BLACKHOLE_THREATS_EXTRA_OPTS:-}"

if [ ! -f "$CONF_PATH" ]; then
    if [ "$CONF_PATH" = "$DEFAULT_CONF" ]; then
        log_info "Creating default configuration file" "config_path=$CONF_PATH"
        cat > "$CONF_PATH" <<'CONF'
# Default configuration
gobgp:
  global:
    config:
      as: 64512
      routerid: "192.168.1.1"
  neighbors:
    - config:
        neighboraddress: "192.168.1.1"
        peeras: 64512

feeds:
#  - url: http://localhost/drop.txt
#    community: 64512:666
CONF
        log_info "Created default configuration file" "config_path=$CONF_PATH"
    else
        log_error "Configuration file not found" "config_path=$CONF_PATH"
        exit 1
    fi
fi

set -- /usr/sbin/blackhole-threats -conf "$CONF_PATH"
if [ -n "$EXTRA_OPTS" ]; then
    # shellcheck disable=SC2086
    set -- "$@" $EXTRA_OPTS
fi

exec s6-setuidgid root "$@"
